# calibre-updater.cron
0 3 * * *  root  /usr/local/bin/calibre-update.sh >>/var/log/calibre-update.log 2>&1
